# etc/sv/exim/run
#!/usr/bin/env bash
exec /usr/sbin/exim -bd -q15m
