# etc/sv/mysql/run
#!/usr/bin/env bash
exec /usr/bin/mysqld_safe --user=mysql 2>&1 >/dev/null
